# deploy/vast/orthrus-budgetmatch.conf
[program:orthrus_budgetmatch]
command=/opt/supervisor-scripts/orthrus-budgetmatch.sh
directory=/workspace/FlowDraft-cfmproper
autostart=false
autorestart=false
startsecs=3
stopwaitsecs=180
stdout_logfile=/dev/shm/flowdraft_runs/orthrus_budgetmatch_20000_r2/supervisor.log
stderr_logfile=/dev/shm/flowdraft_runs/orthrus_budgetmatch_20000_r2/supervisor.err.log
redirect_stderr=false
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1"
