# scripts/vast_hydraflow_probe_supervisor.conf
[program:hydraflow_probe]
command=/opt/supervisor-scripts/hydraflow-probe.sh
directory=/workspace/FlowDraft-hydraflow
autostart=false
autorestart=false
startsecs=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
