# scripts/vast_eagleflow_parallel_continue_supervisor.conf
[program:eagleflow_parallel_continue]
command=/opt/supervisor-scripts/eagleflow-parallel-continue.sh
directory=/workspace/FlowDraft-cfmproper
autostart=false
autorestart=false
startsecs=3
stopwaitsecs=180
stdout_logfile=/workspace/flowdraft_runs/eagleflow_parallel_continue_20000_r1/supervisor.log
stderr_logfile=/workspace/flowdraft_runs/eagleflow_parallel_continue_20000_r1/supervisor.err.log
redirect_stderr=false
environment=PYTHONUNBUFFERED="1"
