# scripts/vast_hydraflow_advanced_token_supervisor.conf
[program:hydraflow_advanced_token]
command=/opt/supervisor-scripts/hydraflow-advanced-token.sh
directory=/workspace/FlowDraft-cfmproper
autostart=false
autorestart=false
startsecs=3
stdout_logfile=/workspace/flowdraft_runs/hydraflow_advanced_token_5000_r1/supervisor.log
stderr_logfile=/workspace/flowdraft_runs/hydraflow_advanced_token_5000_r1/supervisor.err.log
redirect_stderr=false
environment=PYTHONUNBUFFERED="1"
