# scripts/vast_flowdraft_v3_supervisor.conf
[program:flowdraft_v3]
command=/opt/supervisor-scripts/flowdraft-v3.sh
directory=/workspace/FlowDraft
autostart=false
autorestart=false
startsecs=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
