# scripts/vast_simplex_flow_probe_supervisor.conf
[program:simplex_flow_probe]
command=/opt/supervisor-scripts/simplex-flow-probe.sh
directory=/workspace/FlowDraft-cfmproper
autostart=false
autorestart=false
startsecs=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
