# scripts/vast_hydraflow_advanced_token_kl_supervisor.conf
[program:hydraflow_advanced_token_kl]
command=/opt/supervisor-scripts/hydraflow-advanced-token-kl.sh
directory=/workspace/FlowDraft-cfmproper
autostart=false
autorestart=false
startsecs=3
stdout_logfile=/workspace/flowdraft_runs/hydraflow_advanced_token_kl_5000_r3/supervisor.log
stderr_logfile=/workspace/flowdraft_runs/hydraflow_advanced_token_kl_5000_r3/supervisor.err.log
redirect_stderr=false
environment=PYTHONUNBUFFERED="1"
