# scripts/vast_eagleflow_supervisor.conf
[program:eagleflow_attention_screen]
command=/opt/supervisor-scripts/eagleflow-attention-screen.sh
directory=/workspace/FlowDraft-cfmproper
autostart=false
autorestart=false
startsecs=3
stopwaitsecs=180
stdout_logfile=/workspace/flowdraft_runs/eagleflow_attention_screen_300_r1/supervisor.log
stderr_logfile=/workspace/flowdraft_runs/eagleflow_attention_screen_300_r1/supervisor.err.log
redirect_stderr=false
environment=PYTHONUNBUFFERED="1"
